# util/etc/cron.daily/iftd
#!/bin/sh

if ! [ -f /var/run/iftd.pid ]; then
   rm -f /var/lock/iftd.pid
   /etc/init.d/iftd restart
fi
